# files/etc/crontabs/root
# Compact zram every 12 hours
0 */12 * * * /bin/zram compact
